# hopper.path
[Unit]
Description=Watch for new mail to Claude

[Path]
DirectoryNotEmpty=/var/mail/mail/new

[Install]
WantedBy=multi-user.target
