# etc/hopper.path
[Unit]
Description=Watch for mail

[Path]
DirectoryNotEmpty=/var/mail/mail/new

[Install]
WantedBy=multi-user.target
